# r2.service.example
# Example systemd --user unit for running R2 always-on.
# Copy to ~/.config/systemd/user/r2.service, edit the paths, then:
#   systemctl --user daemon-reload
#   systemctl --user enable --now r2.service
#   loginctl enable-linger "$USER"   # keep it running after logout / across reboots
[Unit]
Description=R2 (AreToo) Discord bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Edit these to match your checkout and Node install:
WorkingDirectory=%h/r2
ExecStart=%h/.nvm/versions/node/v24.18.0/bin/node %h/r2/node_modules/tsx/dist/cli.mjs src/index.ts
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
